# launchd/com.agentic.scan-trigger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agentic.scan-trigger</string>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>date; curl -sf --max-time 5 -X POST http://127.0.0.1:8080/trigger; echo; sleep 45; date; curl -sf --max-time 5 -X POST http://127.0.0.1:8080/trigger; echo</string>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/g4ndr1k/agentic-ai/logs/scan-trigger.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/g4ndr1k/agentic-ai/logs/scan-trigger-err.log</string>
</dict>
</plist>
